# ansible/roles/derigenda/templates/docker-hook.service.j2
[Unit]
Description=Docker Hook

[Service]
Restart=always
ExecStart=/usr/local/bin/docker-hook -t {{ docker_hook_token }} -c systemctl restart docker.latin

[Install]
WantedBy=multi-user.target
