# ansible/roles/derigenda/files/systemd/docker.prune.timer
[Unit]
Description=Run docker.prune.service every 10 minutes

[Timer]
OnCalendar=*:30

[Install]
WantedBy=default.target
