# ansible/roles/derigenda/files/systemd/latin.restart.timer
[Unit]
Description=Weekly restart of latin.service

[Timer]
OnCalendar=Sun 6:0

[Install]
WantedBy=default.target
