# ansible/roles/derigenda/files/systemd/latin.restart.service
[Unit]
Description=Weekly restart

[Service]
Type=oneshot
ExecStart=/usr/bin/systemctl restart docker.latin
