# ansible/roles/derigenda/files/systemd/vacuum.service
[Unit]
Description=Cleans up log files

[Service]
Type=oneshot
ExecStart=/usr/bin/journalctl --vacuum-size=2G
